# packaging/launchd/tech.s16e.hort.daemon.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{.Label}}</string>

    <!--
        Rendered by `hort daemon install` at install-time. Do not edit by hand
        unless you know what you're doing — `hort daemon install` will overwrite
        this file. Re-run the command after upgrading the binary.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>{{.BinaryPath}}</string>
        <string>daemon</string>
        <string>start</string>
    </array>

    <!-- Start when the user logs in and keep it alive. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restart loops if the daemon dies in quick succession. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
{{- if .Home}}
        <key>HOME</key>
        <string>{{.Home}}</string>
{{- end}}
    </dict>

{{- if .Home}}
    <key>WorkingDirectory</key>
    <string>{{.Home}}</string>
{{- end}}

    <!-- Log streams for debugging. Rotate manually when they get large. -->
    <key>StandardOutPath</key>
    <string>{{.StdoutLog}}</string>
    <key>StandardErrorPath</key>
    <string>{{.StderrLog}}</string>

    <!--
        Keep the daemon out of macOS's App Nap so it responds quickly even
        when the terminal isn't foreground. It's a negligible CPU consumer
        when idle (blocks on accept()).
    -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
